# systemd/user/onedrive.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/onedrive --monitor
StandardOutput=null
StandardError=null
